# docker/demo/app_demo/app_server_name.ini
[supervisord]
nodaemon=true

[program:server]
command=bash -c "PYTHONIOENCODING=utf-8 python3 run.py"
stopasgroup=true
killasgroup=true
directory=/apps/app_server_name
stdout_logfile=/apps/app_server_name/server_std.log
stderr_logfile=/apps/app_server_name/server_err.log
autostart=true
autorestart=true
priority=1


