# docker/demo/app_demo/celery_service.ini
[supervisord]
nodaemon=true


[program:celery_default_task]
command=bash -c "celery -A app.celery_app worker -l INFO -f ./logs/celery_default_task.log -c 3"
stopasgroup=true
killasgroup=true
directory=/apps/app_server_name
stdout_logfile=/apps/app_server_name/celery_default_task_std.log
stderr_logfile=/apps/app_server_name/celery_default_task_err.log
autostart=true
autorestart=true
priority=2


[program:special_queue_task]
command=bash -c "celery -A app.celery_app worker -l INFO -f ./logs/special_queue_task.log -c 1 -Q queue_name"
stopasgroup=true
killasgroup=true
directory=/apps/app_server_name
stdout_logfile=/apps/app_server_name/special_queue_task_std.log
stderr_logfile=/apps/app_server_name/special_queue_task_err.log
autostart=true
autorestart=true
priority=2


[program:celery_beat]
command=bash -c "celery -A app.celery_app beat -l INFO -f ./logs/celery_beat.log"
stopasgroup=true
killasgroup=true
directory=/apps/app_server_name
stdout_logfile=/apps/app_server_name/celery_beat_std.log
stderr_logfile=/apps/app_server_name/celery_beat_err.log
autostart=true
autorestart=true
priority=2

